# systemd/stored.service
[Unit]
Description=Stored Daemon
Documentation=https://github.com/artur-augustyniak/stored

[Service]
Type=notify
ExecStart=/usr/sbin/stored -f /etc/stored/stored.cfg

[Install]
WantedBy=multi-user.target